# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/ElasticpotPY/
command=/usr/bin/python3 main.py
autorestart=true
